# scripts/snarky.rc
: <<'EOF'
Snarky - Zero Knowledge Dead Drop
Copyright (C) 2026 Sapadian LLC.

This program is free software: you can redistribute it and/or modify
it under the terms of the GNU Affero General Public License as published
by the Free Software Foundation, either version 3 of the License, or
(at your option) any later version.

This program is distributed in the hope that it will be useful,
but WITHOUT ANY WARRANTY; without even the implied warranty of
MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
GNU Affero General Public License for more details.

You should have received a copy of the GNU Affero General Public License
along with this program.  If not, see https://www.gnu.org/licenses/
EOF

#!/bin/sh

# PROVIDE: snarky
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="snarky"
rcvar="snarky_enable"

load_rc_config $name

# CHANGED: Renamed variable to avoid rc.subr auto-switching user
: ${snarky_enable:="NO"}
: ${snarky_daemon_user:="nobody"}
: ${snarky_daemon_group:="nobody"}
: ${snarky_config:="/usr/local/etc/snarky.json"}

pidfile="/var/run/snarky/${name}.pid"
procname="/usr/sbin/daemon"
command="/usr/sbin/daemon"

start_precmd="snarky_precmd"

snarky_precmd()
{
    # Ensure directory exists with correct ownership
    if [ ! -d "/var/run/snarky" ]; then
        install -d -m 0755 -o ${snarky_daemon_user} -g ${snarky_daemon_group} /var/run/snarky
    fi
}

# -f: redirect I/O (prevents "no screen" crash)
# -u: daemon drops privileges to this user
command_args="-P ${pidfile} -r -f -u ${snarky_daemon_user} ${procname} server -config ${snarky_config}"

run_rc_command "$1"